# etc/supervisor/conf.d/docker.conf
[program:docker]
command=/usr/local/bin/docker
autorestart=true
redirect_stderr=true
stdout_logfile=syslog